# site-cookbooks/wekan-nginx/templates/default/wekan.crontab
12 3 * * * root openssl rand 48 > /etc/letsencrypt/live/<%= @fqdn %>/ticket.key
