# site-cookbooks/nginx/templates/default/ssl.crontab
2 2 25 */1 * root /home/webadm/bin/ssl_renewal.sh
